# config/systemd/tide-web.service
[Unit]
Description=Tide Gateway Web Dashboard
Documentation=https://github.com/bodegga/tide
After=network.target tor.service nginx.service
Wants=tor.service
Requires=nginx.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/tide
ExecStart=/usr/bin/python3 /opt/tide/scripts/runtime/tide-web-dashboard.py
Restart=always
RestartSec=10
# ZERO-LOG POLICY: No logging for privacy (Tide is a privacy appliance)
# ABSOLUTELY NO LOGS - not even errors (privacy over debugging)
StandardOutput=null
StandardError=null

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Port 8080 (internal) - nginx proxies port 80 externally
# No special capabilities needed for non-privileged port

[Install]
WantedBy=multi-user.target
